# provision/roles/Stouts.iptables/templates/etc/network/if-up.d/iptables_load.j2
#!/bin/bash

# This file was generated by Ae for {{ansible_fqdn}}
# Do NOT modify this file by hand!

{{iptables_rules_path}}
exit $?
